# services/.config/systemd/user/check-updates.timer
[Unit]
Description=Trigger update check every day at 18h

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=false

[Install]
WantedBy=default.target