# services/.config/systemd/user/check-updates.timer
[Unit]
Description=Trigger update check every day at 18h

[Timer]
OnBootSec=1h
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=default.target