# services/.config/systemd/user/check-updates.service
[Unit]
Description=Check for system updates

[Service]
Type=oneshot
ExecStart=/home/diogo/.dotfiles/scripts/.local/bin/check-updates
